# examples/com.agentic-memorizer.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service label - must match filename -->
    <key>Label</key>
    <string>com.agentic-memorizer.daemon</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/.local/bin/agentic-memorizer</string>
        <string>daemon</string>
        <string>start</string>
    </array>

    <!-- Run at load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep daemon alive - restart if it exits -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Standard output and error logs -->
    <key>StandardOutPath</key>
    <string>/tmp/agentic-memorizer.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agentic-memorizer.err.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <!-- Set HOME for config file resolution -->
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>

        <!-- Optional: Set Claude API key if not in config file -->
        <!-- <key>ANTHROPIC_API_KEY</key> -->
        <!-- <string>your-api-key-here</string> -->
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME</string>

    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Throttle restart attempts to prevent rapid restart loops -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
